# 7.4/mysql/worker.ini
[program:worker]
command = php /var/www/artisan queue:work
user = www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/laravel-worker.log
startretries=50